# system/files/etc/systemd/system/polkit-agent-helper@.service.d/override.conf
[Service]
ReadWritePaths=/var/lib/faillock
